# whisper.cpp/rootfs/etc/services.d/wyoming/run
#!/usr/bin/with-contenv bashio

cd /api
./script/run --uri tcp://0.0.0.0:10303 --debug --api http://localhost:8080/inference